# fixtures/dummy.service
[Unit]
Description=dummy unit for e2e tests

[Service]
ExecStart=/bin/sleep 400
